# ops/daily-signal-emma.service
[Unit]
Description=Daily Signal editorial update on OpenClaw
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=ubuntu
SupplementaryGroups=docker
WorkingDirectory=/opt/openclaw/data/workspace/daily-signal
Environment=PYTHONUNBUFFERED=1
Environment="GIT_SSH_COMMAND=ssh -i /home/ubuntu/.ssh/daily-signal_deploy -o IdentitiesOnly=yes"
EnvironmentFile=-/etc/default/daily-signal-emma
ExecStart=/opt/openclaw/data/workspace/daily-signal/ops/run-emma-update.sh
TimeoutStartSec=90min
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
