# ops/daily-signal-emma.timer
[Unit]
Description=Run the AI digest every morning

[Timer]
OnCalendar=*-*-* 07:00:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=2m
Unit=daily-signal-emma.service

[Install]
WantedBy=timers.target
